# luxtuner.service
[Unit]
Description=Lux Tuner Service
After=network.target

[Service]
Type=simple
User=100acresheater
WorkingDirectory=/home/100acresheater
ExecStart=/home/100acresheater/lux.sh
Environment="PYTHONPATH=$PYTHONPATH:/home/100acresheater/"
Restart=on-failure

[Install]
WantedBy=multi-user.target
